# config/rules.yaml
# ThreatWeave — Windows MITRE ATT&CK Detection Rules
# Each rule maps a Windows Event ID (and optional conditions) to a technique.
# Conditions are ANDed. Supported keys:
#   event_id         : exact match on EventID string
#   process_contains : substring in process_name OR raw_message (case-insensitive)
#   message_contains : substring in raw_message (case-insensitive)
#   username_contains: substring in username field (case-insensitive)

rules:

  # ── Initial Access ──────────────────────────────────────────────────────
  - id: T1078-LOGON-SUCCESS
    technique_id: T1078
    technique_name: Valid Accounts — Successful Logon
    tactic: initial-access
    confidence: medium
    conditions:
      event_id: "4624"

  - id: T1078-LOGON-EXPLICIT
    technique_id: T1078
    technique_name: Valid Accounts — Explicit Credentials
    tactic: initial-access
    confidence: medium
    conditions:
      event_id: "4648"

  # ── Credential Access ───────────────────────────────────────────────────
  - id: T1110-LOGON-FAIL
    technique_id: T1110
    technique_name: Brute Force — Failed Logon
    tactic: credential-access
    confidence: medium
    conditions:
      event_id: "4625"

  - id: T1003-MIMIKATZ
    technique_id: T1003.001
    technique_name: OS Credential Dumping — LSASS (mimikatz)
    tactic: credential-access
    confidence: high
    conditions:
      event_id: "4688"
      process_contains: "mimikatz"

  - id: T1003-PROCDUMP
    technique_id: T1003.001
    technique_name: OS Credential Dumping — LSASS (procdump)
    tactic: credential-access
    confidence: high
    conditions:
      event_id: "4688"
      process_contains: "procdump"

  - id: T1003-SAM
    technique_id: T1003.002
    technique_name: OS Credential Dumping — SAM (reg save)
    tactic: credential-access
    confidence: high
    conditions:
      event_id: "4688"
      process_contains: "reg save"

  - id: T1558-KERBEROS
    technique_id: T1558.003
    technique_name: Steal or Forge Kerberos Tickets — Kerberoasting
    tactic: credential-access
    confidence: medium
    conditions:
      event_id: "4769"

  - id: T1552-FINDSTR
    technique_id: T1552.001
    technique_name: Credentials in Files — findstr
    tactic: credential-access
    confidence: low
    conditions:
      event_id: "4688"
      process_contains: "findstr"

  # ── Execution ───────────────────────────────────────────────────────────
  - id: T1059-POWERSHELL
    technique_id: T1059.001
    technique_name: Command Scripting — PowerShell
    tactic: execution
    confidence: high
    conditions:
      event_id: "4688"
      process_contains: "powershell"

  - id: T1059-POWERSHELL-ENCODED
    technique_id: T1059.001
    technique_name: Command Scripting — PowerShell Encoded
    tactic: execution
    confidence: high
    conditions:
      event_id: "4688"
      message_contains: "encodedCommand"

  - id: T1059-CMD
    technique_id: T1059.003
    technique_name: Command Scripting — Windows CMD
    tactic: execution
    confidence: medium
    conditions:
      event_id: "4688"
      process_contains: "cmd.exe"

  - id: T1059-WSCRIPT
    technique_id: T1059.005
    technique_name: Command Scripting — VBScript (wscript)
    tactic: execution
    confidence: high
    conditions:
      event_id: "4688"
      process_contains: "wscript"

  - id: T1059-CSCRIPT
    technique_id: T1059.005
    technique_name: Command Scripting — VBScript (cscript)
    tactic: execution
    confidence: high
    conditions:
      event_id: "4688"
      process_contains: "cscript"

  - id: T1059-PYTHON
    technique_id: T1059.006
    technique_name: Command Scripting — Python
    tactic: execution
    confidence: medium
    conditions:
      event_id: "4688"
      process_contains: "python"

  - id: T1218-RUNDLL32
    technique_id: T1218.011
    technique_name: Signed Binary Proxy — Rundll32
    tactic: defense-evasion
    confidence: medium
    conditions:
      event_id: "4688"
      process_contains: "rundll32"

  - id: T1218-REGSVR32
    technique_id: T1218.010
    technique_name: Signed Binary Proxy — Regsvr32
    tactic: defense-evasion
    confidence: medium
    conditions:
      event_id: "4688"
      process_contains: "regsvr32"

  - id: T1218-MSHTA
    technique_id: T1218.005
    technique_name: Signed Binary Proxy — Mshta
    tactic: defense-evasion
    confidence: high
    conditions:
      event_id: "4688"
      process_contains: "mshta"

  # ── Persistence ─────────────────────────────────────────────────────────
  - id: T1053-TASK-CREATED
    technique_id: T1053.005
    technique_name: Scheduled Task — Created
    tactic: persistence
    confidence: high
    conditions:
      event_id: "4698"

  - id: T1053-TASK-DELETED
    technique_id: T1053.005
    technique_name: Scheduled Task — Deleted (cleanup)
    tactic: defense-evasion
    confidence: medium
    conditions:
      event_id: "4699"

  - id: T1543-SERVICE
    technique_id: T1543.003
    technique_name: Create or Modify System Process — Windows Service
    tactic: persistence
    confidence: high
    conditions:
      event_id: "7045"

  - id: T1547-REGISTRY
    technique_id: T1547.001
    technique_name: Boot/Logon Autostart — Registry Run Key
    tactic: persistence
    confidence: medium
    conditions:
      event_id: "4657"
      message_contains: "Run"

  - id: T1098-USER-CREATED
    technique_id: T1098
    technique_name: Account Manipulation — User Created
    tactic: persistence
    confidence: high
    conditions:
      event_id: "4720"

  - id: T1098-GROUP-ADD
    technique_id: T1098
    technique_name: Account Manipulation — User Added to Group
    tactic: persistence
    confidence: high
    conditions:
      event_id: "4732"

  # ── Privilege Escalation ────────────────────────────────────────────────
  - id: T1068-SPECIAL-PRIV
    technique_id: T1068
    technique_name: Exploitation for Privilege Escalation — Special Privileges
    tactic: privilege-escalation
    confidence: medium
    conditions:
      event_id: "4672"

  - id: T1548-RUNAS
    technique_id: T1548.002
    technique_name: Abuse Elevation Control — RunAs
    tactic: privilege-escalation
    confidence: medium
    conditions:
      event_id: "4688"
      process_contains: "runas"

  # ── Defense Evasion ─────────────────────────────────────────────────────
  - id: T1070-LOG-CLEAR
    technique_id: T1070.001
    technique_name: Indicator Removal — Clear Windows Event Logs
    tactic: defense-evasion
    confidence: high
    conditions:
      event_id: "1102"

  - id: T1562-AUDIT-POLICY
    technique_id: T1562.001
    technique_name: Impair Defenses — Modify Audit Policy
    tactic: defense-evasion
    confidence: high
    conditions:
      event_id: "4719"

  - id: T1070-VSSADMIN
    technique_id: T1070.004
    technique_name: Indicator Removal — Shadow Copy Deletion (vssadmin)
    tactic: defense-evasion
    confidence: high
    conditions:
      event_id: "4688"
      process_contains: "vssadmin"

  - id: T1070-BCDEDIT
    technique_id: T1490
    technique_name: Inhibit System Recovery — bcdedit
    tactic: impact
    confidence: high
    conditions:
      event_id: "4688"
      process_contains: "bcdedit"

  - id: T1027-OBFUSCATE
    technique_id: T1027
    technique_name: Obfuscated Files or Information — PowerShell Encoding
    tactic: defense-evasion
    confidence: high
    conditions:
      event_id: "4688"
      message_contains: "-enc"

  # ── Discovery ───────────────────────────────────────────────────────────
  - id: T1082-SYSINFO
    technique_id: T1082
    technique_name: System Information Discovery — systeminfo
    tactic: discovery
    confidence: medium
    conditions:
      event_id: "4688"
      process_contains: "systeminfo"

  - id: T1016-IPCONFIG
    technique_id: T1016
    technique_name: System Network Config Discovery — ipconfig
    tactic: discovery
    confidence: medium
    conditions:
      event_id: "4688"
      process_contains: "ipconfig"

  - id: T1018-NET-VIEW
    technique_id: T1018
    technique_name: Remote System Discovery — net view
    tactic: discovery
    confidence: medium
    conditions:
      event_id: "4688"
      message_contains: "net view"

  - id: T1087-NET-USER
    technique_id: T1087.001
    technique_name: Account Discovery — net user
    tactic: discovery
    confidence: medium
    conditions:
      event_id: "4688"
      message_contains: "net user"

  - id: T1135-SHARE
    technique_id: T1135
    technique_name: Network Share Discovery — Share Accessed
    tactic: discovery
    confidence: medium
    conditions:
      event_id: "5140"

  - id: T1083-DIR
    technique_id: T1083
    technique_name: File and Directory Discovery — dir command
    tactic: discovery
    confidence: low
    conditions:
      event_id: "4688"
      process_contains: "dir"

  # ── Lateral Movement ────────────────────────────────────────────────────
  - id: T1021-RDP
    technique_id: T1021.001
    technique_name: Remote Services — Remote Desktop Protocol
    tactic: lateral-movement
    confidence: high
    conditions:
      event_id: "4624"
      message_contains: "LogonType: 10"

  - id: T1021-SMB
    technique_id: T1021.002
    technique_name: Remote Services — SMB / Windows Admin Shares
    tactic: lateral-movement
    confidence: medium
    conditions:
      event_id: "5140"

  - id: T1550-PTH
    technique_id: T1550.002
    technique_name: Use Alternate Auth — Pass the Hash
    tactic: lateral-movement
    confidence: high
    conditions:
      event_id: "4624"
      message_contains: "NtLmSsp"

  # ── Collection ──────────────────────────────────────────────────────────
  - id: T1005-XCOPY
    technique_id: T1005
    technique_name: Data from Local System — xcopy
    tactic: collection
    confidence: low
    conditions:
      event_id: "4688"
      process_contains: "xcopy"

  - id: T1074-7ZIP
    technique_id: T1074.001
    technique_name: Data Staged — 7-Zip archive
    tactic: collection
    confidence: medium
    conditions:
      event_id: "4688"
      process_contains: "7z"

  # ── Impact ──────────────────────────────────────────────────────────────
  - id: T1486-RANSOMWARE
    technique_id: T1486
    technique_name: Data Encrypted for Impact — vssadmin delete
    tactic: impact
    confidence: high
    conditions:
      event_id: "4688"
      message_contains: "vssadmin delete"

  - id: T1489-NET-STOP
    technique_id: T1489
    technique_name: Service Stop — net stop
    tactic: impact
    confidence: medium
    conditions:
      event_id: "4688"
      message_contains: "net stop"
